# autostart/Dropdown Terminal.desktop
[Desktop Entry]
Exec=qterminal --drop
Name=Dropdown Terminal
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
