# autostart/redshift.desktop
[Desktop Entry]
Exec=redshift
Name=redshift
OnlyShowIn=LXQt;
Type=Application
Version=1.0
